# docker/php/crontab
0 0 * * * php /var/www/app/bin/console app:load-rates